# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
pidfile=/var/run/supervisord.pid
user=root

; ─── PHP-FPM ──────────────────────────────────────────────────────────────────
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─── Laravel Scheduler (runs every minute, replaces crontab) ──────────────────
[program:scheduler]
command=/bin/bash -c "while true; do php /var/www/html/artisan schedule:run --no-ansi >> /var/log/scheduler.log 2>&1; sleep 60; done"
autostart=true
autorestart=true
priority=10
user=www-data
stdout_logfile=/var/log/supervisor/scheduler.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/supervisor/scheduler-error.log
stderr_logfile_maxbytes=1MB
